# vnx-autoconf/open-source/unix/upstart/vnxace.conf
description     "VNX Autoconfiguration and Command Execution Daemon (VNXACED)"
start on startup

post-stop script
   kill -INT `cat /var/run/vnxaced.pid`
   rm /var/run/vnxaced.pid
end script


#expect fork
#task
exec `perl /usr/local/bin/vnxaced`