# vnx-autoconf/unix/vnxdaemon.conf
description     "VNX Autoconfiguration daemon"
start on startup
expect fork
exec perl /etc/init.d/vnxdaemon.pl